# nord.timer
[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=10s
Unit=nord.service

[Install]
WantedBy=gnome-session-initialized.target
